# images/s6-busybox/etc/service/busybox/run
#!/bin/sh
echo S6 BusyBox
exec /bin/ash
